# vm-agent.service
[Unit]
Description=VM Agent for Hyperlight Agents
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/vm-agent.sh
Restart=always
RestartSec=5
User=root
StandardOutput=append:/var/log/vm-agent.log
StandardError=append:/var/log/vm-agent.log
Environment=VM_ID=%H

# Security settings
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=false

[Install]
WantedBy=multi-user.target
